# deploy/crynova-worker.conf
; /etc/supervisor/conf.d/crynova-worker.conf
; Очередь: платежи, вебхуки, проверки блокчейна.
; После создания:  supervisorctl reread && supervisorctl update && supervisorctl start crynova-worker:*
[program:crynova-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/crynova/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/crynova
autostart=true
autorestart=true
user=crynova
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/crynova/storage/logs/worker.log
stopwaitsecs=3600
